# ansible/roles/redis/templates/redis-supervisord.conf.j2
[program:redis]
command=/usr/bin/redis-server
user=redis
autorestart=true
stdout_logfile = /var/log/supervisord/%(program_name)s.log
stderr_logfile = /var/log/supervisord/%(program_name)s.error.log
